# pkgs/lumen@.service
[Unit]
Description=Lumen Wayland WebRTC compositor (%i)
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i
PAMName=login
SyslogIdentifier=lumen@%i
Environment=XDG_RUNTIME_DIR=/run/user/%U
Environment=LUMEN_STATIC_DIR=/usr/share/lumen/web
Environment=LUMEN_LOG_OUTPUT=journald
Environment=LUMEN_SYSLOG_IDENTIFIER=lumen@%i
Environment=LUMEN_DESKTOP=labwc
EnvironmentFile=-/etc/lumen/%i.env
EnvironmentFile=-%h/.config/lumen/env
ExecStart=/usr/bin/lumen
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
