# pkgs/debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Create /etc/lumen/ config directory if it doesn't exist
        if [ ! -d /etc/lumen ]; then
            mkdir -p /etc/lumen
            chmod 755 /etc/lumen
        fi

        # Install example config only if not already present
        if [ ! -f /etc/lumen/example.env ]; then
            cp /usr/share/lumen/example.env /etc/lumen/example.env
            chmod 644 /etc/lumen/example.env
        fi

        echo ""
        echo "╔══════════════════════════════════════════════════════════╗"
        echo "║              Lumen — Getting Started                     ║"
        echo "╠══════════════════════════════════════════════════════════╣"
        echo "║  1. Create a config file for the user you want to run:  ║"
        echo "║     System-wide (any user):                              ║"
        echo "║     cp /etc/lumen/example.env /etc/lumen/<username>.env  ║"
        echo "║     nano /etc/lumen/<username>.env                       ║"
        echo "║                                                          ║"
        echo "║     Or per-user (overrides system config):               ║"
        echo "║     mkdir -p ~/.config/lumen                             ║"
        echo "║     cp /etc/lumen/example.env ~/.config/lumen/env        ║"
        echo "║     nano ~/.config/lumen/env                             ║"
        echo "║                                                          ║"
        echo "║  2. Start the service:                                   ║"
        echo "║     systemctl start lumen@<username>                     ║"
        echo "║                                                          ║"
        echo "║  3. Enable autostart on boot (optional):                 ║"
        echo "║     systemctl enable lumen@<username>                    ║"
        echo "║                                                          ║"
        echo "║  4. View logs:                                           ║"
        echo "║     journalctl -u lumen@<username> -f                    ║"
        echo "╚══════════════════════════════════════════════════════════╝"
        echo ""

        # Reload udev rules so the uinput controller rule takes effect immediately
        if command -v udevadm > /dev/null 2>&1; then
            udevadm control --reload-rules && udevadm trigger || true
        fi
        ;;
esac

exit 0
